# hosts/h01/crontab
*/15	*	*	*	*	root	/usr/libexec/hprop h02
